# onsun.service
[Unit]
Description=Schedule sun-based events

[Service]
Type=oneshot
ExecStart=/etc/onsun/onsun.py

[Install]
WantedBy=timers.target
